# systemd/twitcast-promo.timer
[Unit]
Description=TWiT Transcript Promo timer

[Timer]
OnCalendar=*-*-* *:00,30:00
Persistent=true

[Install]
WantedBy=timers.target
