# systemd/twitcast-dashboard.timer
[Unit]
Description=TWiT E-Ink Dashboard timer

[Timer]
OnCalendar=*-*-* 11:00:00 US/Pacific
Persistent=true

[Install]
WantedBy=timers.target
